# junkmap.upstart
# JunkMap Service
start on runlevel [2345]
stop on runlevel [016]

setuid postfix
setgid postfix

respawn
exec /usr/local/bin/junkmap --database /var/lib/postfix/junkmap.db --log /var/log/junkmap.log --domain example.org --address root


